# etc/rc.d/rc.local
#!/bin/sh
#
# /etc/rc.d/rc.local:  Local system initialization script.
#
# Put any local startup commands in here.  Also, if you have
# anything that needs to be run at shutdown time you can
# make an /etc/rc.d/rc.local_shutdown script and put those
# commands in there.

# Mouse didn't work in 3.9-rc2 unless the module was reinserted.
rmmod chromeos_laptop
modprobe chromeos_laptop

# Makes first suspend fail instead of crashing.
modprobe tpm_tis force=1 interrupts=0
